# scripts/jarvis-update.service
[Unit]
Description=JARVIS Auto-Updater (single run)

[Service]
Type=oneshot
User=root
ExecStart=/home/pi/jarvis-voice-assistant-master/scripts/jarvis-update.sh
