# scripts/jarvis-update.service
[Unit]
Description=JARVIS Auto-Updater (single run)

[Service]
Type=oneshot
User=root
ExecStart=__REPO__/scripts/jarvis-update.sh
